# build/install/deb/debian/onlyofficeAppserverStudioNotify.service
[Unit]
Description=ONLYOFFICE AppServer Studio Notify Service

[Service]
PIDFile=/run/appserver/appserver-studio_notify.pid
Type=simple
User=onlyoffice
Group=onlyoffice
WorkingDirectory=/var/www/services/studio.notify/
ExecStart=ASC.Studio.Notify.dll
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target