# build/install/deb/debian/onlyofficeAppserverBackup.service
[Unit]
Description=ONLYOFFICE AppServer Backup Service

[Service]
PIDFile=/run/appserver/appserver-backup.pid
Type=simple
User=onlyoffice
Group=onlyoffice
WorkingDirectory=/var/www/services/backup/
ExecStart=ASC.Data.Backup.dll
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target