# build/install/deb/debian/onlyofficeAppserverPeopleServer.service
[Unit]
Description=ONLYOFFICE AppServer People Server Service

[Service]
PIDFile=/run/appserver/appserver-people_server.pid
Type=simple
User=onlyoffice
Group=onlyoffice
WorkingDirectory=/var/www/products/ASC.People/server/
ExecStart=ASC.People.dll
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target