# build/install/deb/debian/onlyofficeAppserverFilesServices.service
[Unit]
Description=ONLYOFFICE AppServer Files Services Service

[Service]
PIDFile=/run/appserver/appserver-files_services.pid
Type=simple
User=onlyoffice
Group=onlyoffice
WorkingDirectory=/var/www/products/ASC.Files/service/
ExecStart=ASC.Files.Service.dll
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target