# build/install/deb/debian/onlyofficeAppserverApi.service
[Unit]
Description=ONLYOFFICE AppServer API Service

[Service]
PIDFile=/run/appserver/appserver-api.pid
Type=simple
User=onlyoffice
Group=onlyoffice
WorkingDirectory=/var/www/studio/api/
ExecStart=ASC.Web.Api.dll
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target