# build/install/deb/debian/onlyofficeAppserverStudio.service
[Unit]
Description=ONLYOFFICE AppServer Studio Service

[Service]
PIDFile=/run/appserver/appserver-studio.pid
Type=simple
User=onlyoffice
Group=onlyoffice
WorkingDirectory=/var/www/studio/server/
ExecStart=ASC.Web.Studio.dll
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target