# build/install/deb/debian/onlyofficeAppserverUrlshortener.service
[Unit]
Description=ONLYOFFICE AppServer Urlshortener Service

[Service]
PIDFile=/run/appserver/appserver-urlshortener.pid
Type=simple
User=onlyoffice
Group=onlyoffice
WorkingDirectory=/services/urlshortener/service/
ExecStart=ASC.UrlShortener.Svc.dll
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target