# build/install/deb/debian/onlyofficeAppserverSocket.service
[Unit]
Description=ONLYOFFICE AppServer Socket Service

[Service]
PIDFile=/run/appserver/appserver-socket.pid
Type=simple
User=onlyoffice
Group=onlyoffice
WorkingDirectory=/services/socket/service/
ExecStart=ASC.Socket.IO.Svc.dll
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target