# docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:ASC.Data.Backup.BackgroundTasks]
command=dotnet ASC.Data.Backup.BackgroundTasks.dll --urls=http://0.0.0.0:5032 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Data.Backup.BackgroundTasks core:eventBus:subscriptionClientName=asc_event_bus_backup_queue
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Data.Backup.BackgroundTasks/service
autorestart=true
priority=300

[program:ASC.ApiSystem]
command=dotnet ASC.ApiSystem.dll --urls=http://0.0.0.0:5010 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.ApiSystem
directory=%(ENV_SRC_PATH)s/publish/services/ASC.ApiSystem/service
autorestart=true
priority=300

[program:ASC.ClearEvents]
command=dotnet ASC.ClearEvents.dll --urls=http://0.0.0.0:5027 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.ClearEvents
directory=%(ENV_SRC_PATH)s/publish/services/ASC.ClearEvents/service
autorestart=true
priority=300

[program:ASC.Data.Backup]
command=dotnet ASC.Data.Backup.dll --urls=http://0.0.0.0:5012 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Data.Backup
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Data.Backup/service
autorestart=true
priority=300

[program:ASC.Files]
command=dotnet ASC.Files.dll --urls=http://0.0.0.0:5007 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Files
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Files/service
autorestart=true
priority=300

[program:ASC.Files.Service]
command=dotnet ASC.Files.Service.dll --urls=http://0.0.0.0:5009 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Files.Service
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Files.Service/service
autorestart=true
priority=300

[program:ASC.Notify]
command=dotnet ASC.Notify.dll --urls=http://0.0.0.0:5005 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Notify core:eventBus:subscriptionClientName=asc_event_bus_notify_queue
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Notify/service
autorestart=true
priority=300

[program:ASC.People]
command=dotnet ASC.People.dll --urls=http://0.0.0.0:5004 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.People
directory=%(ENV_SRC_PATH)s/publish/services/ASC.People/service
autorestart=true
priority=300

[program:ASC.Studio.Notify]
command=dotnet ASC.Studio.Notify.dll --urls=http://0.0.0.0:5006 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Studio.Notify
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Studio.Notify/service
autorestart=true
priority=300

[program:ASC.Web.Api]
command=dotnet ASC.Web.Api.dll --urls=http://0.0.0.0:5000 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Web.Api
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Web.Api/service
autorestart=true
priority=300

[program:ASC.Web.Studio]
command=dotnet ASC.Web.Studio.dll --urls=http://0.0.0.0:5003 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Web.Studio core:eventBus:subscriptionClientName=asc_event_bus_webstudio_queue
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Web.Studio/service
autorestart=true
priority=300

[program:ASC.Web.HealthChecks.UI]
command=dotnet ASC.Web.HealthChecks.UI.dll --urls=http://0.0.0.0:5033 --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_ENV_EXTENSION)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Web.HealthChecks.UI
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Web.HealthChecks.UI/service
autorestart=true
priority=300

[program:ASC.Sdk]
command=node server.js --app.port=5099 --app.appsettings=%(ENV_PATH_TO_CONF)s   
directory=%(ENV_SRC_PATH)s/publish/web/sdk
autorestart=true
priority=300

[program:ASC.Editors]
command=node server.js --app.port=5013 --app.appsettings=%(ENV_PATH_TO_CONF)s   
directory=%(ENV_SRC_PATH)s/publish/web/editor
autorestart=true
priority=300

[program:ASC.Login]
command=node server.js --app.port=5011 --app.appsettings=%(ENV_PATH_TO_CONF)s   
directory=%(ENV_SRC_PATH)s/publish/web/login
autorestart=true
priority=300

[program:ASC.Socket.IO]
command=node server.js --app.port=9899 --app.appsettings=%(ENV_PATH_TO_CONF)s   
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Socket.IO
autorestart=true
priority=300

[program:ASC.SsoAuth]
command=node app.js --app.port=9834 --app.appsettings=%(ENV_PATH_TO_CONF)s   
directory=%(ENV_SRC_PATH)s/publish/services/ASC.SsoAuth
autorestart=true
priority=300

[program:ASC.Identity.Authorization]
command=java -jar app.jar
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Identity.Authorization
autorestart=true
stderr_logfile=/var/log/supervisor/ASC.Identity.Authorization.err.log
stdout_logfile=/var/log/supervisor/ASC.Identity.Authorization.out.log
environment=
    JAVA_HOME="/usr/lib/jvm/java-21-openjdk-amd64",
    SERVER_PORT=8080,
    LOG_FILE_PATH="/var/log/onlyoffice/ASC.Identity.Authorization.log"
priority=300

[program:ASC.Identity.Registration]
command=java -jar app.jar
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Identity.Registration
autorestart=true
stderr_logfile=/var/log/supervisor/ASC.Identity.Registration.err.log
stdout_logfile=/var/log/supervisor/ASC.Identity.Registration.out.log
environment=
    JAVA_HOME="/usr/lib/jvm/java-21-openjdk-amd64",
    SERVER_PORT=9090,
    LOG_FILE_PATH="/var/log/onlyoffice/ASC.Identity.Registration.log"
priority=300

[program:init_nginx]
command=/nginx/docker-entrypoint.sh
autorestart=false
startretries=1
priority=200

[program:openresty]
command=bash -c 'sleep 20 && /usr/local/openresty/bin/openresty -g "daemon off;"'
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/openresty_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/openresty_stdout.log
stdout_logfile_maxbytes=10MB
priority=350

[program:redis]
command=/usr/bin/redis-server
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/redis.err.log
stdout_logfile=/var/log/supervisor/redis.out.log
priority=200

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/rabbitmq.err.log
stdout_logfile=/var/log/supervisor/rabbitmq.out.log
priority=200

[program:mysql]
command=/usr/sbin/mysqld
user=mysql
autostart=true
autorestart=true
environment=HOME=/var/lib/mysql
stderr_logfile=/var/log/mysql/error.log
stdout_logfile=/var/log/mysql/out.log
priority=100

[program:init_db]
command=bash -c 'sleep 5 && /init_db.sh' 
user=root
autostart=true
autorestart=false
startretries=2
startsecs=10
priority=150

[program:ASC.Migration.Runner]
command=bash -c 'echo "MigArg: %(ENV_MIGRATION_TYPE)s" &&  sleep 15 && dotnet ASC.Migration.Runner.dll "%(ENV_MIGRATION_TYPE)s"'
directory=%(ENV_SRC_PATH)s/publish/services/ASC.Migration.Runner/service
autostart=true
autorestart=false
startretries=3
startsecs=20
priority=250

[program:opensearch]
command=/opt/opensearch/bin/opensearch
user=opensearch
directory=/opt/opensearch
autostart=true
autorestart=true
startsecs=30
stopwaitsecs=90
stopsignal=INT
killasgroup=true
environment=JAVA_HOME="/opt/opensearch/jdk", OPENSEARCH_JAVA_OPTS="-Xms1g -Xmx1g -Dopensearch.path.home=/opt/opensearch"
stdout_logfile=/var/log/opensearch/opensearch.out.log
stderr_logfile=/var/log/opensearch/opensearch.err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
priority=250

[program:opensearch-dashboards]
command=/opt/opensearch-dashboards/bin/opensearch-dashboards
user=opensearch
directory=/opt/opensearch-dashboards
autostart=true
autorestart=true
startsecs=30
stopwaitsecs=90
stopsignal=INT
killasgroup=true
environment=NODE_OPTIONS="--max-old-space-size=1024"
stdout_logfile=/var/log/opensearch/dashboards.out.log
stderr_logfile=/var/log/opensearch/dashboards.err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
priority=250