# build/install/docker/config/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/usr/bin/pidproxy /var/mysqld/mysqld.pid /usr/bin/mysqld_safe --pid-file=/var/mysqld/mysqld.pid
autostart=true
autorestart=true
environment=PATH=/usr/local/sbin:/usr/bin:/bin:/usr/local/bin
user=mysql

[program:kafka]
directory=/root/kafka_2.12-2.5.0/
command=/root/kafka_2.12-2.5.0/bin/kafka-server-start.sh /root/kafka_2.12-2.5.0/config/server.properties
autostart=true
autorestart=true

[program:kafka_zookeeper]
directory=/root/kafka_2.12-2.5.0/
command=/root/kafka_2.12-2.5.0/bin/zookeeper-server-start.sh /root/kafka_2.12-2.5.0/config/zookeeper.properties
autostart=true
autorestart=true

[program:api]
directory=/var/www/studio/api/
command=dotnet ASC.Web.Api.dll --urls=http://0.0.0.0:5000 --pathToConf=/app/onlyoffice/config/ --STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=api --ENVIRONMENT=test
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:api_system]
directory=/var/www/services/apisystem
command=dotnet ASC.ApiSystem.dll --urls=http://0.0.0.0:5010 --pathToConf=/app/onlyoffice/config/ --STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=apisystem --ENVIRONMENT=test
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:urlshortener]
directory=/services/urlshortener/service 
command=dotnet ASC.UrlShortener.Svc.dll --urls=http://0.0.0.0:5015 --pathToConf=/app/onlyoffice/config/ --STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=urlshortener --ENVIRONMENT=test
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:thumbnails]
directory=/services/thumbnails/service 
command=dotnet ASC.Thumbnails.Svc.dll --urls=http://0.0.0.0:5016 --pathToConf=/app/onlyoffice/config/ --STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=thumbnails 
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:studio_notify]
directory=/var/www/services/studio.notify/
command=dotnet ASC.Studio.Notify.dll --urls=http://0.0.0.0:5006 --pathToConf=/app/onlyoffice/config/ --STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=notify.studio  --core:products:folder=/var/www/products/ --core:products:subfolder=server --ENVIRONMENT=test
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:notify]
directory=/var/www/services/notify/
command=dotnet ASC.Notify.dll --urls=http://0.0.0.0:5005 --pathToConf=/app/onlyoffice/config/ --STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=notify  --core:products:folder=/var/www/products/ --core:products:subfolder=server
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:people]
directory=/var/www/products/ASC.People/server/
command=dotnet ASC.People.dll --urls=http://0.0.0.0:5004 --pathToConf=/app/onlyoffice/config/  --STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice/ --log:name=people --ENVIRONMENT=test
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:files]
directory=/var/www/products/ASC.Files/server/
command=dotnet ASC.Files.dll --urls=http://0.0.0.0:5007 --pathToConf=/app/onlyoffice/config/  --STORAGE_ROOT=/app/onlyoffice/data/  --log:dir=/var/log/onlyoffice/ --log:name=files --ENVIRONMENT=test
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:files_services]
directory=/var/www/products/ASC.Files/service
command=dotnet ASC.Files.Service.dll --urls=http://0.0.0.0:5009 --pathToConf=/app/onlyoffice/config/  --STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice/ --log:name=filesservice  --core:products:folder=/var/www/products --core:products:subfolder=server --ENVIRONMENT=test
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:studio]
directory=/var/www/studio/server/
command=dotnet ASC.Web.Studio.dll --urls=http://0.0.0.0:5003 --pathToConf=/app/onlyoffice/config/  --STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=web --ENVIRONMENT=test
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:backup]
directory=/var/www/studio/backup/
command=dotnet ASC.Data.Backup.dll --urls=http://0.0.0.0:5012 --pathToConf=/app/onlyoffice/config/  --STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=backup --core:products:folder=/var/www/products/ --core:products:subfolder=server
autostart=true
autorestart=true
user=root
environment=ASPNETCORE_ENVIRONMENT=test

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
